# data/rm-pad.service
[Unit]
Description=rm-pad input forwarding for reMarkable tablet

[Service]
Type=simple
ExecStart=/usr/bin/rm-pad
Restart=always
RestartSec=5
Environment="RMPAD_HOST=10.11.99.1"

[Install]
WantedBy=default.target
